# deploy/layrr-server.service
[Unit]
Description=Layrr Server (Hono process manager)
After=network.target

[Service]
Type=simple
User=layrr
WorkingDirectory=/opt/layrr/packages/server
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
EnvironmentFile=/opt/layrr/.env

[Install]
WantedBy=multi-user.target
